# dot_config/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland session target
BindsTo=graphical-session.target
Wants=xdg-desktop-autostart.target
