# dot_config/systemd/user/tldr-cache-update.service
[Unit]
Description=Updates `tldr` cache

[Service]
Type=oneshot
ExecStart=/usr/bin/tldr --update

[Install]
WantedBy=multi-user.target
